# toolsets/falcon-py27-nginx-gunicorn-meinheld/supervisord.conf
[supervisord]

[unix_http_server]
file=/tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock


[program:app_server]
environment=PORT="8001"
command=foreman start -p 8001
priority=1
startretries=1
autorestart=false
redirect_stderr=true
stopasgroup=true ; stop subprocesses (started by foreman)
stdout_logfile=/tmp/app_server.log
stdout_logfile_maxbytes=100KB
stdout_logfile_backups=0

[program:nginx]
environment=PORT="8000"
command=../../../nginx -p . -c ../../../nginx_proxy.conf
priority=999
startretries=1
autorestart=false
stopasgroup=true ; stop subprocesses
